# supervisor/horizon-tasks-marinsan-scool-cat_LOCAL.conf
[program:horizon-tasks-marinsan-scool-cat]
process_name=%(program_name)s
command=php /home/marin/Code/Marinsan/tasks/artisan horizon
autostart=true
autorestart=true
user=marin
redirect_stderr=true
stdout_logfile=/home/marin/Code/Marinsan/tasks/storage/logs/horizon.log
